# config/supervisor/7-haproxy.conf
[program:HAProxy]
command=haproxy -f /config/haproxy.cfg
priority=4
autorestart=true
stderr_logfile=/var/log/haproxy.err.log
stdout_logfile=/var/log/haproxy.out.log